# cf-p-link.service
[Unit]
Description=Cloudflare Persistent Tunnel Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cf-p-link
ExecStart=/usr/local/bin/tunnel-daemon-native
Restart=always
RestartSec=10
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin"
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
